# var/lib/pacman/local/man-db-2.9.2-1/files
%FILES%
etc/
etc/man_db.conf
usr/
usr/bin/
usr/bin/accessdb.exe
usr/bin/apropos.exe
usr/bin/catman.exe
usr/bin/convert-mans
usr/bin/lexgrog.exe
usr/bin/man-recode.exe
usr/bin/man.exe
usr/bin/mandb.exe
usr/bin/manpath.exe
usr/bin/whatis.exe
usr/lib/
usr/lib/man-db/
usr/lib/man-db/globbing.exe
usr/lib/man-db/libman.a
usr/lib/man-db/libmandb.a
usr/lib/man-db/manconv.exe
usr/lib/man-db/zsoelim.exe
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/man-db.service
usr/lib/systemd/system/man-db.timer
usr/share/
usr/share/doc/
usr/share/doc/man-db/
usr/share/doc/man-db/man-db-manual.ps
usr/share/doc/man-db/man-db-manual.txt
usr/share/locale/
usr/share/locale/af/
usr/share/locale/af/LC_MESSAGES/
usr/share/locale/af/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ast/
usr/share/locale/ast/LC_MESSAGES/
usr/share/locale/ast/LC_MESSAGES/man-db.mo
usr/share/locale/be/
usr/share/locale/be/LC_MESSAGES/
usr/share/locale/be/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ca/LC_MESSAGES/man-db.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/cs/LC_MESSAGES/man-db.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/da/LC_MESSAGES/man-db.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/de/LC_MESSAGES/man-db.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/eo/
usr/share/locale/eo/LC_MESSAGES/
usr/share/locale/eo/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/eo/LC_MESSAGES/man-db.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/es/LC_MESSAGES/man-db.mo
usr/share/locale/et/
usr/share/locale/et/LC_MESSAGES/
usr/share/locale/et/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/fi/LC_MESSAGES/man-db.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/fr/LC_MESSAGES/man-db.mo
usr/share/locale/ga/
usr/share/locale/ga/LC_MESSAGES/
usr/share/locale/ga/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/man-db.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/it/LC_MESSAGES/man-db.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ja/LC_MESSAGES/man-db.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ms/
usr/share/locale/ms/LC_MESSAGES/
usr/share/locale/ms/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/nl/LC_MESSAGES/man-db.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/pl/LC_MESSAGES/man-db.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/pt/LC_MESSAGES/man-db.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/pt_BR/LC_MESSAGES/man-db.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ro/LC_MESSAGES/man-db.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/ru/LC_MESSAGES/man-db.mo
usr/share/locale/rw/
usr/share/locale/rw/LC_MESSAGES/
usr/share/locale/rw/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/sr/LC_MESSAGES/man-db.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/sv/LC_MESSAGES/man-db.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/tr/LC_MESSAGES/man-db.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/vi/
usr/share/locale/vi/LC_MESSAGES/
usr/share/locale/vi/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/vi/LC_MESSAGES/man-db.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/zh_CN/LC_MESSAGES/man-db.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/man-db-gnulib.mo
usr/share/locale/zh_TW/LC_MESSAGES/man-db.mo
usr/share/man/
usr/share/man/da/
usr/share/man/da/man1/
usr/share/man/da/man1/apropos.1.gz
usr/share/man/da/man1/lexgrog.1.gz
usr/share/man/da/man1/man-recode.1.gz
usr/share/man/da/man1/man.1.gz
usr/share/man/da/man1/manconv.1.gz
usr/share/man/da/man1/manpath.1.gz
usr/share/man/da/man1/whatis.1.gz
usr/share/man/da/man1/zsoelim.1.gz
usr/share/man/da/man5/
usr/share/man/da/man5/manpath.5.gz
usr/share/man/da/man8/
usr/share/man/da/man8/accessdb.8.gz
usr/share/man/da/man8/catman.8.gz
usr/share/man/da/man8/mandb.8.gz
usr/share/man/de/
usr/share/man/de/man1/
usr/share/man/de/man1/apropos.1.gz
usr/share/man/de/man1/lexgrog.1.gz
usr/share/man/de/man1/man-recode.1.gz
usr/share/man/de/man1/man.1.gz
usr/share/man/de/man1/manconv.1.gz
usr/share/man/de/man1/manpath.1.gz
usr/share/man/de/man1/whatis.1.gz
usr/share/man/de/man1/zsoelim.1.gz
usr/share/man/de/man5/
usr/share/man/de/man5/manpath.5.gz
usr/share/man/de/man8/
usr/share/man/de/man8/accessdb.8.gz
usr/share/man/de/man8/catman.8.gz
usr/share/man/de/man8/mandb.8.gz
usr/share/man/es/
usr/share/man/es/man1/
usr/share/man/es/man1/apropos.1.gz
usr/share/man/es/man1/lexgrog.1.gz
usr/share/man/es/man1/man-recode.1.gz
usr/share/man/es/man1/man.1.gz
usr/share/man/es/man1/manconv.1.gz
usr/share/man/es/man1/manpath.1.gz
usr/share/man/es/man1/whatis.1.gz
usr/share/man/es/man1/zsoelim.1.gz
usr/share/man/es/man5/
usr/share/man/es/man5/manpath.5.gz
usr/share/man/es/man8/
usr/share/man/es/man8/accessdb.8.gz
usr/share/man/es/man8/catman.8.gz
usr/share/man/es/man8/mandb.8.gz
usr/share/man/fr/
usr/share/man/fr/man1/
usr/share/man/fr/man1/apropos.1.gz
usr/share/man/fr/man1/lexgrog.1.gz
usr/share/man/fr/man1/man-recode.1.gz
usr/share/man/fr/man1/man.1.gz
usr/share/man/fr/man1/manconv.1.gz
usr/share/man/fr/man1/manpath.1.gz
usr/share/man/fr/man1/whatis.1.gz
usr/share/man/fr/man1/zsoelim.1.gz
usr/share/man/fr/man5/
usr/share/man/fr/man5/manpath.5.gz
usr/share/man/fr/man8/
usr/share/man/fr/man8/accessdb.8.gz
usr/share/man/fr/man8/catman.8.gz
usr/share/man/fr/man8/mandb.8.gz
usr/share/man/id/
usr/share/man/id/man1/
usr/share/man/id/man1/apropos.1.gz
usr/share/man/id/man1/lexgrog.1.gz
usr/share/man/id/man1/man-recode.1.gz
usr/share/man/id/man1/man.1.gz
usr/share/man/id/man1/manconv.1.gz
usr/share/man/id/man1/manpath.1.gz
usr/share/man/id/man1/whatis.1.gz
usr/share/man/id/man1/zsoelim.1.gz
usr/share/man/id/man5/
usr/share/man/id/man5/manpath.5.gz
usr/share/man/id/man8/
usr/share/man/id/man8/accessdb.8.gz
usr/share/man/id/man8/catman.8.gz
usr/share/man/id/man8/mandb.8.gz
usr/share/man/it/
usr/share/man/it/man1/
usr/share/man/it/man1/apropos.1.gz
usr/share/man/it/man1/man.1.gz
usr/share/man/it/man1/manpath.1.gz
usr/share/man/it/man1/whatis.1.gz
usr/share/man/it/man1/zsoelim.1.gz
usr/share/man/it/man5/
usr/share/man/it/man5/manpath.5.gz
usr/share/man/it/man8/
usr/share/man/it/man8/accessdb.8.gz
usr/share/man/it/man8/catman.8.gz
usr/share/man/it/man8/mandb.8.gz
usr/share/man/ja/
usr/share/man/ja/man1/
usr/share/man/ja/man1/apropos.1.gz
usr/share/man/ja/man1/lexgrog.1.gz
usr/share/man/ja/man1/man-recode.1.gz
usr/share/man/ja/man1/man.1.gz
usr/share/man/ja/man1/manconv.1.gz
usr/share/man/ja/man1/manpath.1.gz
usr/share/man/ja/man1/whatis.1.gz
usr/share/man/ja/man1/zsoelim.1.gz
usr/share/man/ja/man5/
usr/share/man/ja/man5/manpath.5.gz
usr/share/man/ja/man8/
usr/share/man/ja/man8/accessdb.8.gz
usr/share/man/ja/man8/catman.8.gz
usr/share/man/ja/man8/mandb.8.gz
usr/share/man/man1/
usr/share/man/man1/apropos.1.gz
usr/share/man/man1/lexgrog.1.gz
usr/share/man/man1/man-recode.1.gz
usr/share/man/man1/man.1.gz
usr/share/man/man1/manconv.1.gz
usr/share/man/man1/manpath.1.gz
usr/share/man/man1/whatis.1.gz
usr/share/man/man1/zsoelim.1.gz
usr/share/man/man5/
usr/share/man/man5/manpath.5.gz
usr/share/man/man8/
usr/share/man/man8/accessdb.8.gz
usr/share/man/man8/catman.8.gz
usr/share/man/man8/mandb.8.gz
usr/share/man/nl/
usr/share/man/nl/man1/
usr/share/man/nl/man1/apropos.1.gz
usr/share/man/nl/man1/lexgrog.1.gz
usr/share/man/nl/man1/man-recode.1.gz
usr/share/man/nl/man1/man.1.gz
usr/share/man/nl/man1/manconv.1.gz
usr/share/man/nl/man1/manpath.1.gz
usr/share/man/nl/man1/whatis.1.gz
usr/share/man/nl/man1/zsoelim.1.gz
usr/share/man/nl/man5/
usr/share/man/nl/man5/manpath.5.gz
usr/share/man/nl/man8/
usr/share/man/nl/man8/accessdb.8.gz
usr/share/man/nl/man8/catman.8.gz
usr/share/man/nl/man8/mandb.8.gz
usr/share/man/pl/
usr/share/man/pl/man1/
usr/share/man/pl/man1/apropos.1.gz
usr/share/man/pl/man1/lexgrog.1.gz
usr/share/man/pl/man1/man-recode.1.gz
usr/share/man/pl/man1/man.1.gz
usr/share/man/pl/man1/manconv.1.gz
usr/share/man/pl/man1/manpath.1.gz
usr/share/man/pl/man1/whatis.1.gz
usr/share/man/pl/man1/zsoelim.1.gz
usr/share/man/pl/man5/
usr/share/man/pl/man5/manpath.5.gz
usr/share/man/pl/man8/
usr/share/man/pl/man8/accessdb.8.gz
usr/share/man/pl/man8/catman.8.gz
usr/share/man/pl/man8/mandb.8.gz
usr/share/man/pt/
usr/share/man/pt/man1/
usr/share/man/pt/man1/apropos.1.gz
usr/share/man/pt/man1/lexgrog.1.gz
usr/share/man/pt/man1/man-recode.1.gz
usr/share/man/pt/man1/man.1.gz
usr/share/man/pt/man1/manconv.1.gz
usr/share/man/pt/man1/manpath.1.gz
usr/share/man/pt/man1/whatis.1.gz
usr/share/man/pt/man1/zsoelim.1.gz
usr/share/man/pt/man5/
usr/share/man/pt/man5/manpath.5.gz
usr/share/man/pt/man8/
usr/share/man/pt/man8/accessdb.8.gz
usr/share/man/pt/man8/catman.8.gz
usr/share/man/pt/man8/mandb.8.gz
usr/share/man/pt_BR/
usr/share/man/pt_BR/man1/
usr/share/man/pt_BR/man1/apropos.1.gz
usr/share/man/pt_BR/man1/lexgrog.1.gz
usr/share/man/pt_BR/man1/man-recode.1.gz
usr/share/man/pt_BR/man1/man.1.gz
usr/share/man/pt_BR/man1/manconv.1.gz
usr/share/man/pt_BR/man1/manpath.1.gz
usr/share/man/pt_BR/man1/whatis.1.gz
usr/share/man/pt_BR/man1/zsoelim.1.gz
usr/share/man/pt_BR/man5/
usr/share/man/pt_BR/man5/manpath.5.gz
usr/share/man/pt_BR/man8/
usr/share/man/pt_BR/man8/accessdb.8.gz
usr/share/man/pt_BR/man8/catman.8.gz
usr/share/man/pt_BR/man8/mandb.8.gz
usr/share/man/ru/
usr/share/man/ru/man1/
usr/share/man/ru/man1/apropos.1.gz
usr/share/man/ru/man1/lexgrog.1.gz
usr/share/man/ru/man1/man-recode.1.gz
usr/share/man/ru/man1/man.1.gz
usr/share/man/ru/man1/manconv.1.gz
usr/share/man/ru/man1/manpath.1.gz
usr/share/man/ru/man1/whatis.1.gz
usr/share/man/ru/man1/zsoelim.1.gz
usr/share/man/ru/man5/
usr/share/man/ru/man5/manpath.5.gz
usr/share/man/ru/man8/
usr/share/man/ru/man8/accessdb.8.gz
usr/share/man/ru/man8/catman.8.gz
usr/share/man/ru/man8/mandb.8.gz
usr/share/man/sr/
usr/share/man/sr/man1/
usr/share/man/sr/man1/apropos.1.gz
usr/share/man/sr/man1/lexgrog.1.gz
usr/share/man/sr/man1/man-recode.1.gz
usr/share/man/sr/man1/man.1.gz
usr/share/man/sr/man1/manconv.1.gz
usr/share/man/sr/man1/manpath.1.gz
usr/share/man/sr/man1/whatis.1.gz
usr/share/man/sr/man1/zsoelim.1.gz
usr/share/man/sr/man5/
usr/share/man/sr/man5/manpath.5.gz
usr/share/man/sr/man8/
usr/share/man/sr/man8/accessdb.8.gz
usr/share/man/sr/man8/catman.8.gz
usr/share/man/sr/man8/mandb.8.gz
usr/share/man/sv/
usr/share/man/sv/man1/
usr/share/man/sv/man1/apropos.1.gz
usr/share/man/sv/man1/lexgrog.1.gz
usr/share/man/sv/man1/man-recode.1.gz
usr/share/man/sv/man1/man.1.gz
usr/share/man/sv/man1/manconv.1.gz
usr/share/man/sv/man1/manpath.1.gz
usr/share/man/sv/man1/whatis.1.gz
usr/share/man/sv/man1/zsoelim.1.gz
usr/share/man/sv/man5/
usr/share/man/sv/man5/manpath.5.gz
usr/share/man/sv/man8/
usr/share/man/sv/man8/accessdb.8.gz
usr/share/man/sv/man8/catman.8.gz
usr/share/man/sv/man8/mandb.8.gz
usr/share/man/tr/
usr/share/man/tr/man1/
usr/share/man/tr/man1/apropos.1.gz
usr/share/man/tr/man1/lexgrog.1.gz
usr/share/man/tr/man1/man-recode.1.gz
usr/share/man/tr/man1/man.1.gz
usr/share/man/tr/man1/manconv.1.gz
usr/share/man/tr/man1/manpath.1.gz
usr/share/man/tr/man1/whatis.1.gz
usr/share/man/tr/man1/zsoelim.1.gz
usr/share/man/tr/man5/
usr/share/man/tr/man5/manpath.5.gz
usr/share/man/tr/man8/
usr/share/man/tr/man8/accessdb.8.gz
usr/share/man/tr/man8/catman.8.gz
usr/share/man/tr/man8/mandb.8.gz
usr/share/man/zh_CN/
usr/share/man/zh_CN/man1/
usr/share/man/zh_CN/man1/apropos.1.gz
usr/share/man/zh_CN/man1/lexgrog.1.gz
usr/share/man/zh_CN/man1/man-recode.1.gz
usr/share/man/zh_CN/man1/man.1.gz
usr/share/man/zh_CN/man1/manconv.1.gz
usr/share/man/zh_CN/man1/manpath.1.gz
usr/share/man/zh_CN/man1/whatis.1.gz
usr/share/man/zh_CN/man1/zsoelim.1.gz
usr/share/man/zh_CN/man5/
usr/share/man/zh_CN/man5/manpath.5.gz
usr/share/man/zh_CN/man8/
usr/share/man/zh_CN/man8/accessdb.8.gz
usr/share/man/zh_CN/man8/catman.8.gz
usr/share/man/zh_CN/man8/mandb.8.gz

%BACKUP%
etc/man_db.conf	f791dd86450fabf98b60b7c13143c0d8

